# salt/openstack/horizon/horizon.conf
description "OpenStack Horozon"
author "ChenHao Ku <maple52046@gmail.com>"

start on runlevel [2345]
stop on runlevel [016]

respawn

exec uwsgi --ini /etc/sscloud/horizon.ini
